# crontab
# crontab
# m h dom mon dow command
*/10 * * * * root /usr/local/bin/python /home/jirapullapp/jirapull.py &> /var/log/cron.log 2>&1
# required blank row
